# supervisord.conf
; AutoApply process supervisor (single-server production).
;
; AutoApply's runtime is three long-lived processes:
;
;   web    -- FastAPI + SPA, the operator console (`autoapply web`)
;   worker -- Celery worker consuming search / materials /
;             application / maintenance queues (`autoapply worker`)
;   beat   -- Celery Beat with redbeat (`autoapply beat`); the
;             single cron driver. Only ONE beat instance should run.
;
; The README's section 14 shows a single systemd unit for the web;
; that is fine for evaluation but skips the worker and beat. This
; supervisord config is the lightest way to run all three with auto-
; restart, log rotation, and dependency-free upgrades, without the
; ceremony of three separate systemd units.
;
; Usage on a Linux server:
;
;   sudo apt install -y supervisor
;   sudo cp /opt/autoapply/supervisord.conf /etc/supervisor/conf.d/autoapply.conf
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl status
;
; To run unprivileged (recommended -- matches the README's
; `autoapply` service user), install the user-mode supervisor and
; point it at this file:
;
;   supervisord -c /opt/autoapply/supervisord.conf
;
; Edit the four paths marked `EDIT:` below before deploying.

[supervisord]
; Local working dir; supervisord drops its own pidfile + log here.
; Override via `-l` / `-j` if you prefer /var/log + /var/run.
logfile=/opt/autoapply/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=5
loglevel=info
pidfile=/opt/autoapply/logs/supervisord.pid
nodaemon=false

[unix_http_server]
; `supervisorctl` talks to supervisord over this socket. 0700 so
; only the service user can issue start/stop/restart.
file=/opt/autoapply/logs/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///opt/autoapply/logs/supervisor.sock

; ---- Shared program defaults ----
;
; All three programs share the same working dir, environment, and
; service user. Edit the four `EDIT:` lines and `directory` once;
; the per-program blocks inherit from these by repetition.

; ---------- Web GUI ----------
[program:autoapply-web]
; `--no-open` because there is no desktop to launch a browser on;
; `--host 127.0.0.1` because the public entry point is Nginx (see
; README section 16). Change the host to 0.0.0.0 only if you are
; binding directly to the public interface.
command=/opt/autoapply/.local/bin/uv run autoapply web --host 127.0.0.1 --port 8000 --no-open --show-logs
directory=/opt/autoapply
; EDIT: service user that owns /opt/autoapply.
user=autoapply
autostart=true
autorestart=true
startsecs=5
; SIGTERM is correct for uvicorn -- it drains in-flight requests.
stopsignal=TERM
stopwaitsecs=30
; supervisord captures stdout/stderr; rotate at 10MB x 5 to keep
; disk usage bounded without losing recent context.
stdout_logfile=/opt/autoapply/logs/web.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/opt/autoapply/logs/web.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
; EDIT: HOME is required so `uv` finds its cache and the LLM CLIs
; (Claude Code / Codex CLI) find their credential stores.
environment=HOME="/opt/autoapply",PATH="/opt/autoapply/.local/bin:/usr/local/bin:/usr/bin:/bin"

; ---------- Celery worker ----------
[program:autoapply-worker]
; Default `--concurrency 2` matches `cmd_worker.py`'s default; AutoApply
; tasks are minutes-long (materials generation, browser fill) so high
; concurrency is wasteful. Tune `-c` per host (see D025 in DECISIONS).
command=/opt/autoapply/.local/bin/uv run autoapply worker
directory=/opt/autoapply
user=autoapply
autostart=true
autorestart=true
startsecs=10
; SIGTERM -> Celery's warm-shutdown: finish in-flight tasks then
; exit. stopwaitsecs is generous because a browser-fill task can
; take 1-2 min to drain cleanly.
stopsignal=TERM
stopwaitsecs=120
stdout_logfile=/opt/autoapply/logs/worker.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/opt/autoapply/logs/worker.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=HOME="/opt/autoapply",PATH="/opt/autoapply/.local/bin:/usr/local/bin:/usr/bin:/bin"

; ---------- Celery Beat (cron driver) ----------
[program:autoapply-beat]
; Beat is the single cron source. redbeat (configured in
; src/tasks/app.py) provides leader election via Redis so running
; this on more than one host is *safe* but unnecessary -- pick one
; box as the beat host.
command=/opt/autoapply/.local/bin/uv run autoapply beat
directory=/opt/autoapply
user=autoapply
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopwaitsecs=30
stdout_logfile=/opt/autoapply/logs/beat.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/opt/autoapply/logs/beat.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=HOME="/opt/autoapply",PATH="/opt/autoapply/.local/bin:/usr/local/bin:/usr/bin:/bin"

; ---- Optional: group control ----
;
; `supervisorctl restart autoapply:` restarts all three at once.
; Useful after `git pull && uv sync && alembic upgrade head`.
[group:autoapply]
programs=autoapply-web,autoapply-worker,autoapply-beat
